# video2text.service
[Unit]
Description=video2text WebUI - 视频转字幕服务
Documentation=https://github.com/your-repo/video2text
After=network.target
Wants=network.target

[Service]
Type=simple
User=lym
Group=lym
WorkingDirectory=/home/lym/projects/video2text

# 启动前清理占用 7881 端口的残留进程，避免 Gradio 报端口被占用
ExecStartPre=/usr/bin/bash -lc '/usr/bin/fuser -k 7881/tcp >/dev/null 2>&1 || true'

# 使用 venv 内的 Python
ExecStart=/home/lym/projects/video2text/.venv/bin/python fastapi_app.py --host 0.0.0.0 --port 7881 --ssl-certfile /home/lym/projects/video2text/video2text.pem --ssl-keyfile /home/lym/projects/video2text/video2text-key.pem

# 崩溃后自动重启，最多 3 次（5 分钟内）
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=3

# 环境变量
Environment="PYTHONUNBUFFERED=1"
Environment="GRADIO_ANALYTICS_ENABLED=false"
# 如需指定 CUDA 设备，取消注释：
# Environment="CUDA_VISIBLE_DEVICES=0"

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video2text

[Install]
WantedBy=multi-user.target
